# linux/install/DEBIAN/postrm
#!/bin/sh
update-rc.d chromeboxcontroller remove
