# linux/install/DEBIAN/postinst
#!/bin/sh
update-rc.d chromeboxcontroller defaults
service chromeboxcontroller start
